# com.arxiv.crawler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.arxiv.crawler</string>
  <key>ProgramArguments</key><array>
    <string>/Users/zhengxinyu/.hanako/skills/arxiv-paper-digest/.venv/bin/python</string>
    <string>/Users/zhengxinyu/.hanako/skills/arxiv-paper-digest/run.py</string>
    <string>--once</string>
  </array>
  <key>StartCalendarInterval</key><dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>0</integer></dict>
  <key>StandardOutPath</key><string>/Users/zhengxinyu/.hanako/skills/arxiv-paper-digest/logs/crawler-scheduled.log</string>
  <key>StandardErrorPath</key><string>/Users/zhengxinyu/.hanako/skills/arxiv-paper-digest/logs/crawler-scheduled-error.log</string>
</dict></plist>
